# minecraft-bedrock.service
[Unit]
Description=Minecraft Bedrock Server
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
# screen detaches from the invoking process, so let systemd know the service
# will fork and manage the resulting screen session.
Type=forking
User=minecraft
Group=minecraft
WorkingDirectory=/opt/minecraft-bedrock
ExecStart=/usr/bin/screen -DmS bedrock /bin/bash -lc './start.sh'
ExecStop=/usr/bin/screen -S bedrock -X quit
Restart=on-failure
RestartSec=5
UMask=0027

# Systemd hardening options
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
CapabilityBoundingSet=
AmbientCapabilities=
ReadWritePaths=/opt/minecraft-bedrock

[Install]
WantedBy=multi-user.target
